# src/main/dist/bin/ponyhub.service
[Unit]
Description=Pony Hub, a search facility for Pony Language developers.

[Service]
User=pony
Group=nogroup
EnvironmentFile=-/etc/default/ponyhub
ExecStart=/usr/local/lib/ponyhub/bin/ponyhub
Restart=on-failure

[Install]
WantedBy=multi-user.target
